# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log

[program:api]
command=/app/main
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/api.log 